# daemon/data/init.d/redhat/nanny
#!/bin/sh
#
# nanny-daemon:    NannyDaemon 
#
# chkconfig: 345 99 50 
# description: Run NannyDaemon 
#
# processname: NannyDaemon
# pidfile: /var/run/nanny/NannyDaemon.pid 
#
